# alykoti/yellow/install/alykoti-web.service
[Unit]
Description=Alykoti Web — Next.js paikallinen käyttöliittymä (Yellow Pi)
After=network-online.target alykoti-yellow.service
Wants=network-online.target
Requires=alykoti-yellow.service

[Service]
Type=simple
User=ek
WorkingDirectory=/home/ek/alykoti-yellow/web
EnvironmentFile=/home/ek/alykoti-yellow/web/.env.local
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
